# curator4-861225780.conf
[program:curator4-861225780]
command=python3 /home/curator4/curator4-861225780/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-861225780/test.err.log
stdout_logfile=/home/curator4/curator4-861225780/test.out.log